# src/bios_alerts_list_server.service.in
[Unit]
Description=bios_alerts_list_server service
After=network.target

[Service]
Type=simple
ExecStart=@prefix@/bin/bios_alerts_list_server @prefix@/etc/@PACKAGE@/bios_alerts_list_server.cfg

[Install]
WantedBy=multi-user.target
