# src/fty-alert-list.service.in
[Unit]
Description=fty-alert-list service
Requires=malamute.target
After=malamute.target
PartOf=bios.target

[Service]
Type=simple
User=nobody
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/sysconfig/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/sysconfig/bios__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-alert-list 

[Install]
WantedBy=bios.target
